# packages/msm-firmware-dumper/files/msm-firmware-dumper.init
#!/bin/sh /etc/rc.common
USE_PROCD=1
START=08
NAME=msm-firmware-dumper

start() {
  [ -f /etc/msm-fw-dumped ] && exit 0
  # /usr/bin/msm-firmware-dumper || true
}
