# packages/msm-firmware-dumper/files/msm-firmware-dumper.init
#!/bin/sh /etc/rc.common
USE_PROCD=1
START=08
NAME=msm-firmware-dumper

export MSM_DUMPER_FLAG_FILE="/etc/msm-fw-dumped"

start() {
  [ -f $MSM_DUMPER_FLAG_FILE ] && exit 0
  /usr/bin/msm-firmware-dumper || true
}
